# .vscode/homeassistant-supervised/etc/systemd/system/systemd-journal-gatewayd.socket.d/10-hassio-supervisor.conf
[Socket]
ListenStream=
ListenStream=/run/systemd-journal-gatewayd.sock
